# systemd/tmuxsaver-restore.service
[Unit]
Description=Restore tmux sessions on login
After=default.target

[Service]
Type=oneshot
# Start a headless tmux server if none is running, then restore sessions.
# The server stays alive in the background; attach with: tmux attach
ExecStart=/bin/bash -c 'tmux start-server; /usr/bin/tmuxsaver restore --quiet'

[Install]
WantedBy=default.target
